# antipenipuan.conf
; Supervisor — fragment program (bukan file utama supervisord.conf)
;
; Pasang:
;   sudo cp /var/www/antipenipuan/antipenipuan.conf /etc/supervisor/conf.d/antipenipuan.conf
;   sudo mkdir -p /tmp/logs
;   sudo chown -R www-data:www-data /tmp/logs
;   Sesuaikan `user`, `command` (path node), dan `directory` di bawah.
;   sudo supervisorctl reread
;   sudo supervisorctl update antipenipuan-form
;   sudo supervisorctl start antipenipuan-form


[program:antipenipuan-form]
directory=/var/www/antipenipuan
command=/usr/bin/node scripts/form-flow-once.mjs
user=www-data
autostart=false
autorestart=true
startsecs=15
startretries=100
stopwaitsecs=120
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/logs/antipenipuan-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/tmp/logs/antipenipuan-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
; Jika perlu override tanpa .env (opsional, hapus titik koma depannya):
;environment=BASE_URL="https://contoh.invalid",LOGIN_PATH="/x/login.html",RUN_COUNT="1"
